# docker/php-fpm/aggregate.cron
*/15 * * * * cd /var/www/pinboard && APP_ENV=${APP_ENV:-dev} APP_DEBUG=${APP_DEBUG:-0} php bin/console aggregate --no-interaction
